# docker/supervisord.conf
[supervisord]
nodaemon=true
user=culifeed
logfile=/app/logs/supervisord.log
pidfile=/app/logs/supervisord.pid
childlogdir=/app/logs

[program:culifeed-bot]
command=python run_bot.py
directory=/app
user=culifeed
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=PYTHONPATH="/app"

[program:culifeed-scheduler]
command=python run_scheduler.py --service
directory=/app
user=culifeed
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=PYTHONPATH="/app"

[unix_http_server]
file=/app/logs/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///app/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface